# deploy/systemd/gigaam-api.service
[Unit]
Description=GigaAM v3 Transcriber API
After=network.target
Documentation=https://github.com/salute-developers/GigaAM

[Service]
Type=simple
User=dubr1k
Group=dubr1k
WorkingDirectory=/mnt/storage10tb/syncthing/development/GigaAMv3

# Переменные окружения
Environment="PATH=/mnt/storage10tb/anaconda/envs/gigaam/bin:/usr/local/bin:/usr/bin:/bin"
EnvironmentFile=/mnt/storage10tb/syncthing/development/GigaAMv3/.env

# Запуск с параметрами из .env
ExecStart=/mnt/storage10tb/anaconda/envs/gigaam/bin/python api.py

# Автоматический перезапуск при сбое
Restart=always
RestartSec=10

# Ограничения ресурсов (опционально)
# LimitNOFILE=65535
# MemoryLimit=8G
# CPUQuota=400%

# Логирование
StandardOutput=journal
StandardError=journal
SyslogIdentifier=gigaam-api

# Безопасность
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target

